# scripts/archiso/airootfs/etc/systemd/system/horizonos-vm-setup.service
[Unit]
Description=HorizonOS VM Guest Tools Setup
After=graphical.target
ConditionPathExists=!/var/lib/horizonos/vm-configured

[Service]
Type=oneshot
ExecStartPre=/bin/mkdir -p /var/lib/horizonos
ExecStart=/usr/local/bin/horizonos-vm-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target