# scripts/archiso/airootfs/etc/systemd/system/horizonos-vm-setup.service
[Unit]
Description=HorizonOS VM Guest Tools Setup
After=network.target
Before=display-manager.service
ConditionPathExists=!/var/lib/horizonos/vm-configured
ConditionVirtualization=yes

[Service]
Type=oneshot
ExecStartPre=/bin/mkdir -p /var/lib/horizonos
ExecStart=/usr/local/bin/horizonos-vm-setup
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target