# scripts/configs/systemd/container-export@.service
[Unit]
Description=HorizonOS Container Binary Export for %i
After=horizonos-containers.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/horizon-container export %i
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target